# .config/services/enabled/autotiling.service
[Unit]
Description=Autotiling

[Service]
ExecStart=autotiling

[Install]
WantedBy=graphical-session.target
