# .config/services/enabled/swayidle.service
[Unit]
Description=swayidle service
Documentation=man:swayidle(1)

[Service]
Type=simple
ExecStart=swayidle
Restart=always
