# .config/services/enabled/syncthing.service
[Unit]
Description=Syncthing

[Service]
ExecStart=syncthing
Restart=on-failure
